# com.glowplan.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.glowplan.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/rameshsubramani/Downloads/TheExperiment/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rameshsubramani/Downloads/TheExperiment</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>3000</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/rameshsubramani/Downloads/TheExperiment/data/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rameshsubramani/Downloads/TheExperiment/data/server.log</string>
</dict>
</plist>
